# addons/backblaze-backup/templates/server-backup.service.j2
[Unit]
Description=Server Backup to Backblaze B2
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart={{ backup_script_path }} backup
User=root
Group=root

Environment="RCLONE_CONFIG={{ rclone_config_file }}"

# Security (ProtectHome=no because rclone.conf is in /root)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths={{ backup_staging_path }} {{ backup_log_path }} /root/.config/rclone
{% for job in backup_jobs %}
ReadOnlyPaths={{ job.source }}
{% endfor %}

# Resource limits
MemoryMax=1G
CPUQuota=80%

StandardOutput=journal
StandardError=journal
SyslogIdentifier=server-backup

TimeoutStartSec=7200
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
